# infra/systemd/deploymagi-admin.service
[Unit]
Description=DeployMagi Admin Dashboard
After=network.target

[Service]
Type=simple
User=ubuntu
Group=ubuntu
WorkingDirectory=/home/ubuntu/deploymagi/admin
ExecStart=/usr/bin/npx serve -s dist -l 4003 --cors
EnvironmentFile=/home/ubuntu/deploymagi/.env
Restart=on-failure
RestartSec=5
SyslogIdentifier=deploymagi-admin
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
