# infra/systemd/numavm-dashboard.service
[Unit]
Description=NumaVM Dashboard
After=network.target

[Service]
Type=simple
User=ubuntu
Group=ubuntu
WorkingDirectory=/home/ubuntu/numavm/dashboard
ExecStart=/usr/bin/npx serve -s dist -l 4002 --cors
EnvironmentFile=/home/ubuntu/numavm/.env
Restart=on-failure
RestartSec=5
SyslogIdentifier=numavm-dashboard
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
