# infra/systemd/numavm-control-plane.service
[Unit]
Description=NumaVM Control Plane
After=network.target numavm-auth.service numavm-bridge.service

[Service]
Type=simple
User=root
Group=root
WorkingDirectory=/home/ubuntu/numavm
ExecStart=/usr/bin/node --import tsx control-plane/server.ts
EnvironmentFile=/home/ubuntu/numavm/.env
Restart=on-failure
RestartSec=5
KillMode=process
SyslogIdentifier=numavm-cp
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
